# root/etc/services.d/syncthing/run
#!/usr/bin/with-contenv bash

if [ -n "$UMASK_SET" ]; then
	umask "$UMASK_SET"
fi

exec \
	s6-setuidgid abc syncthing \
	-home=/config -no-browser -no-restart \
	--gui-address="0.0.0.0:8384"
