# root/etc/services.d/syncthing/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc syncthing \
	-home=/config -no-browser -no-restart \
	--gui-address="0.0.0.0:8384"
